# beaglebone/lib/systemd/system/arch-upgrade.service
[Unit]
Description=Arch Linux automatic upgrade

[Service]
Type=simple
ExecStart=/bin/bash /home/alarm/pacman-upgrade.sh
